# bin/pironman5.service
# https://www.freedesktop.org/software/systemd/man/systemd.service.html
[Unit]
Description=pironman5 service
# Need to start last to avoid gpio being occupied
After=multi-user.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/pironman5
Environment=LD_LIBRARY_PATH=/usr/local/lib
ExecStart=/usr/local/bin/pironman5 start
# PrivateTmp=False

[Install]
WantedBy=multi-user.target
